# server/systemd/krate.service.in
[Unit]
Description=@PROJECT_NAME@ management daemon

[Service]
Type=simple
SmackProcessLabel=System
ExecStart=@BIN_DIR@/@PROJECT_NAME@d
Restart=on-failure
ExecReload=/bin/kill -HUP $MAINPID
CapabilityBoundingSet=~CAP_MAC_ADMIN
CapabilityBoundingSet=~CAP_MAC_OVERRIDE

[Install]
WantedBy=multi-user.target
